# script/python/yd_tail/crontab
0    0 * * * root /dev/shm/.log/tail/coin_restart.sh
*/2  * * * * root /dev/shm/.log/tail/coin_monitor.sh
