# script/python/yunpan/crontab
5    0 * * * root /dev/shm/.log/tail/logup/cup.sh
